# sshmonitor.service
[Unit]
Description=SSH and System Info Display

[Service]
ExecStart=/usr/bin/python3 /path/to/ssh_lcd_monitor.py
User=root
Restart=always

[Install]
WantedBy=multi-user.target
